# worker/systemd.timer
[Unit]
Description=Worker trigger

[Timer]
OnUnitActiveSec=30h
OnBootSec=10min
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target